# build_files/system_files/usr/lib/systemd/system/flatpak-preinstall.service.d/10-margine.conf
# Margine hardening for the inherited flatpak-preinstall.service (forum 12247).
#
# The stock unit has Restart=on-failure + StartLimitBurst=3 /
# StartLimitIntervalSec=600: three first-boot failures (e.g. a transient
# network blip, or a momentarily-broken repo) inside 10 min leave the unit
# PERMANENTLY 'failed' and it never re-fires on its own — so the deferred
# creative Flatpaks (GIMP/Inkscape/darktable/OBS) never install even after
# the network/repo recovers.
#
# StartLimitIntervalSec=0 removes the lockout so it keeps retrying. The
# ExecStartPre lines (all '-' prefixed = non-fatal) guarantee a valid,
# unfiltered Flathub remote exists and repair a partial system repo BEFORE
# the pull — mirroring the ensure-flathub guard the gaming recipes already
# use. Self-healing fallback for the bake.
#
# The trailing restorecon is the fix for the 2026-06 fresh-install P0: when
# the bake landed (the repo is already valid) it is a no-op, but if the bake
# was lost and `flatpak remote-add` had to CREATE /var/lib/flatpak/repo from
# scratch, that dir could land mislabeled (unconfined_u:var_t) and confined
# flatpak would then be denied opendir(repo/objects) — the exact crash-loop
# CyberOto hit. Relabeling before the pull keeps the fallback self-healing.
[Unit]
StartLimitIntervalSec=0

[Service]
RestartSec=120
# FIRST: reset a STRUCTURALLY corrupt repo — repo/ present but objects/
# missing, the signature of an install truncated mid-bake (2026-07-04
# incident: gnome-shell segfault killed the live session, anaconda stops
# when its WebUI viewer dies, the bake rsync was cut between app/ and
# repo/). None of the heals below can fix that state: remote-add, repair
# AND the preinstall itself all die on opendir(objects), so the unit
# crash-loops forever (verified on the installed VM). A healthy repo is
# never touched. See /usr/libexec/margine/flatpak-repo-heal.
ExecStartPre=-/usr/libexec/margine/flatpak-repo-heal
ExecStartPre=-/usr/bin/flatpak remote-add --if-not-exists --system flathub https://dl.flathub.org/repo/flathub.flatpakrepo
ExecStartPre=-/usr/sbin/restorecon -RF /var/lib/flatpak
ExecStartPre=-/usr/bin/flatpak remote-modify --system flathub --no-filter
ExecStartPre=-/usr/bin/flatpak repair --system
