# build_files/system_files/etc/xdg/autostart/easyeffects-service.desktop
[Desktop Entry]
Type=Application
Name=EasyEffects (background service)
Comment=Run the EasyEffects audio DSP chain headless from login
Exec=flatpak run com.github.wwmm.easyeffects --service-mode --hide-window
Icon=com.github.wwmm.easyeffects
NoDisplay=true
X-GNOME-Autostart-enabled=true
# --service-mode runs EasyEffects (Qt, 8.x) as a background service from login;
# --hide-window keeps its window CLOSED at startup. --service-mode ALONE still
# pops the window open in 8.x (user-reported 2026-06-19), so BOTH flags are
# required for a truly hidden start. The DSP chain is live from login but
# applies no effects until a preset is enabled, so this is inert by default;
# re-open EasyEffects from Activities and it attaches to this running service.
# (--gapplication-service is the deprecated 8.x alias of --service-mode.)
#
# Opt out per-user: cp this file to ~/.config/autostart/ with an added
# "Hidden=true" line, or toggle it off in GNOME Settings -> Apps -> Startup.
#
# NOTE: deliberately no X-GNOME-Autostart-Phase (GNOME 50+ SKIPS entries that
# set it; see margine-first-boot.desktop).
