# build_files/system_files/usr/lib/systemd/system/uupd.service.d/10-margine.conf
# Margine: bound uupd's runtime so a hung update step cannot wedge the system,
# and self-heal the known hang class via the watchdog.
#
# uupd.service is Type=oneshot and ships TimeoutStartSec=infinity upstream, so if
# any step hangs the whole run hangs forever. The real-world failure (2026-06-28,
# and AGAIN 2026-07-08): uupd's `flatpak update` got stuck downloading the
# extra-data of fm.reaper.Reaper (the vendor CDN drops the connection, sockets
# sit in CLOSE-WAIT, flatpak's extra-data fetch has no timeout). The hung process
# HELD THE FLATPAK SYSTEM LOCK, so Bazaar and every `flatpak install` blocked
# indefinitely.
#
# A finite TimeoutStartSec makes systemd kill a stuck run (SIGTERM->SIGKILL),
# which releases the flatpak/rpm-ostree locks. uupd.timer (Persistent=true,
# daily 04:00) retries on the next run. rpm-ostree/bootc and flatpak updates are
# transactional/resumable, so a killed run is wasteful at worst, never corrupting.
#
# Two hard-won caveats (2026-07-10 incident analysis):
#  - the timeout counts MONOTONIC time, which freezes across suspend: on a
#    laptop that mostly sleeps, one 30-minute cycle stretched over 38h of
#    calendar time. There is no wall-clock TimeoutStartSec in systemd; the
#    watchdog below is what keeps this bounded IN PRACTICE.
#  - the kill alone does not break the LOOP: upstream Restart=on-failure
#    retries, the retry hangs on the same app (73 cycles observed), and the
#    flatpak lock stays held most of the awake time. OnFailure= hands each
#    timeout to margine-update-watchdog, which strikes the culprit extra-data
#    app and, at 2 strikes in 7 days, `flatpak mask`s it so unattended runs
#    complete again. The user gets a session notification with the recovery
#    command (ujust margine-update-unblock). NOTE: uupd updates the SYSTEM
#    image before flatpaks, so even a machine caught in this loop kept
#    staging OS/security updates: the loop starves apps, not the OS.
[Unit]
OnFailure=margine-update-watchdog.service

[Service]
TimeoutStartSec=30min
