# build_files/system_files/usr/lib/systemd/system/uupd.service.d/10-margine.conf
# Margine: bound uupd's runtime so a hung update step cannot wedge the system.
#
# uupd.service is Type=oneshot and ships TimeoutStartSec=infinity upstream, so if
# any step hangs the whole run hangs forever. The real-world failure (2026-06-28):
# uupd's `flatpak update` got stuck in the `apply_extra` of an extra-data app
# (fm.reaper.Reaper downloads its binary from the vendor and has no internal
# timeout). That hung process HELD THE FLATPAK SYSTEM LOCK for ~1.7 days, so
# Bazaar and every `flatpak install` blocked indefinitely waiting for the lock.
#
# A finite TimeoutStartSec makes systemd kill a stuck run (SIGTERM->SIGKILL),
# which releases the flatpak/rpm-ostree locks. uupd.timer (Persistent=true,
# daily 04:00) retries on the next run. rpm-ostree/bootc and flatpak updates are
# transactional/resumable, so a killed run is wasteful at worst, never corrupting.
#
# 30 min is generous for a normal nightly update (image pull + flatpaks + brew +
# distrobox) on a sane connection, while bounding a hang to minutes instead of days.
[Service]
TimeoutStartSec=30min
