# build_files/system_files/usr/lib/systemd/user/margine-staged-update-notify.timer
[Unit]
Description=Margine: periodic check for long-staged updates

[Timer]
# Late enough that uupd's daily run (and any boot-time staging) has
# settled, then on the wall clock. OnUnitActiveSec alone counted
# MONOTONIC time, which freezes across suspend: on a laptop that
# sleeps most of the day the "every 6h" tick effectively never fired
# (observed 2026-07-10: no tick for 38h of calendar time). OnCalendar
# + Persistent fires promptly after resume instead. The script itself
# rate-limits every notification class.
OnStartupSec=20min
OnCalendar=*-*-* 09,13,17,21:00:00
Persistent=true
AccuracySec=10min

[Install]
WantedBy=timers.target
