# build_files/system_files/usr/lib/systemd/user/margine-staleness.timer
[Unit]
Description=Margine: schedule staleness check every 12h

[Timer]
OnBootSec=10min
OnUnitActiveSec=12h
AccuracySec=10min
# No Persistent=true: it only applies to OnCalendar= timers; with the
# monotonic triggers above it was a silent no-op.

[Install]
WantedBy=timers.target
