# build_files/system_files/usr/lib/systemd/user/margine-upgrade-notify.service
[Unit]
Description=Margine: notify if booted deployment changed since last login
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/margine-upgrade-notify

# graphical-session.target, NOT default.target: default.target comes up
# before the session's notification daemon, so the toast fired into the
# void and (pre-2026-07-10) the script then recorded it as delivered.
[Install]
WantedBy=graphical-session.target
